# deploy/s6-rc.d/sb-embed/run
#!/command/execlineb -P
foreground { mkdir -p /data/logs }
with-contenv
exec sb-stack embed-server
