# deploy/s6-rc.d/sb-mcp/run
#!/command/execlineb -P
with-contenv
exec sb-stack mcp
